# scripts/supervisor.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[group:goo-worker]
priority=999
# programs=nginx,php8-fpm,laravel-schedule,laravel-notification,laravel-queue
programs=go-prog


[program:go-prog]
priority=10
autostart=true
autorestart=true
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
command=/usr/go/bin/main -g 'daemon off;'
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
